# crontab
*/5 * * * * /scripts/crtempfilesync.sh <FORTIMANAGER_USER> <FORTIMANAGER_PASSWORD> <FORTIMANAGER_IP>
0 1 * * * /scripts/crtempfilelic.sh <FORTIMANAGER_USER> <FORTIMANAGER_PASSWORD> <FORTIMANAGER_IP>
