# crontab
*/5 * * * * /scripts/crtempfilesync.sh
0 1 * * * /scripts/crtempfilelic.sh